# afk-guard.service
[Unit]
Description=AFK Guard - Anti-kick input injector
After=graphical-session.target

[Service]
ExecStart=/usr/local/bin/afk-guard
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
